# cattleshed/cattleshed.service.in
[Unit]
Description=Backend of Sandbox Compiler

[Service]
Type=simple
User=cattleshed
Group=cattleshed
WorkingDirectory=@CMAKE_INSTALL_PREFIX@
Restart=on-failure
ExecStart=@CATTLESHED_BINDIR@/cattleshed -c @CATTLESHED_SYSCONFDIR@/cattleshed.conf @CATTLESHED_SERVICE_ARGS@ --syslog /var/log/cattleshed.log

[Install]
WantedBy=multi-user.target
